# class08/node_exporter.ini
[program:node_exporter]
command=/usr/local/bin/node_exporter
stopsignal=QUIT
autostart=true
autorestart=true
startretries=10
startsecs=0
stopwaitsecs=10
stopasgroup=true

stdout_logfile=/var/log/node_exporter/stdout_exporter.log
stdout_logfile_maxbytes=60MB
stdout_logfile_backups=4
stdout_capture_maxbytes=4MB

stderr_logfile=/var/log/node_exporter/stderr_exporter.log
stderr_logfile_maxbytes=60MB
stderr_logfile_backups=4
stderr_capture_maxbytes=4MB